# steps/data/setup/puppet/modules/env/files/xen/xen/xen-g5k.service
[Unit]
Description=Generate MAC addresse for Xen DomU and create /var/log/xen

[Service]
ExecStart=/usr/sbin/xen-g5k start

[Install]
WantedBy=multi-user.target
